# systemd/simple-hyper-trail-worker.timer
[Unit]
Description=Run Simple-Hyper trail/grid worker every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s
Persistent=false
Unit=simple-hyper-trail-worker.service

[Install]
WantedBy=timers.target
